# rootfs/skeleton/etc/inittab

# Unconditional interactive shell on serial console
ttyLP0::askfirst:/bin/sh -l
